# service/rsyslog/run
#!/bin/bash
set -e

# "-n" avoid's auto-backgrounding of the process.
exec rsyslogd -n
